# deployment/docker/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:imgdmp_server]
directory=/app
command=yarn start

[program:imgdmp_housekeeping]
directory=/app
command=yarn housekeeping